# roles/restic/files/restic-backup.service
[Unit]
Description=run the restic backup script

[Service]
ExecStart=/usr/local/bin/restic-backup.sh
