# roles/pg_backup/files/backup.service
[Unit]
Description=Backup the postgresql instance

[Service]
EnvironmentFile=/var/lib/backup/env
ExecStart=/usr/local/bin/backup.sh
