# roles/pg_backup/files/backup.service
[Unit]
Description=Backup the postgresql instance

[Service]
ExecStart=/usr/local/bin/backup.sh
